# packaging/systemd/zeno-agent.service
[Unit]
Description=Zeno agent
After=network-online.target{{AGENT_AFTER_CONTROLLER}}
Wants=network-online.target
{{AGENT_REQUIRES_CONTROLLER}}
[Service]
Type=simple
User={{RUN_USER}}
WorkingDirectory={{AGENT_WORKING_DIR}}
ExecStart={{AGENT_BINARY}} -controller-url {{CONTROLLER_URL}} -node-id {{NODE_ID}} -token-file {{AGENT_TOKEN_FILE}} -state-interval {{AGENT_STATE_INTERVAL}} -heartbeat-interval {{AGENT_HEARTBEAT_INTERVAL}} -host-interval {{AGENT_HOST_INTERVAL}} -identity-refresh-interval {{AGENT_IDENTITY_REFRESH_INTERVAL}} -version {{AGENT_VERSION}}{{AGENT_EXTRA_ARGS}}
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
