# packaging/systemd/zeno-controller.service
[Unit]
Description=Zeno controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RUN_USER}}
WorkingDirectory={{INSTALL_DIR}}/current
ExecStart={{CONTROLLER_BINARY}} -addr {{CONTROLLER_ADDR}} -web-dir {{INSTALL_DIR}}/current/web -db {{DATA_DIR}}/zeno.db{{SEED_PREVIEW_FLAG}} -agent-token-file {{AGENT_TOKEN_FILE}} -admin-token-file {{ADMIN_TOKEN_FILE}}{{AGENT_BINARY_FLAG}}{{AGENT_VERSION_FLAG}}
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths={{INSTALL_DIR}}
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
